# templates/etc/profile.d/java.sh.j2
# {{ ansible_managed }}
export JAVA_HOME="{{ _java_home[java_release] }}"
